# setup/claude-telegram.service
[Unit]
Description=Claude Code Telegram Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
User=claudeuser
WorkingDirectory=/home/claudeuser
Environment=HOME=/home/claudeuser
Environment=TERM=xterm-256color
Environment=PATH=/home/claudeuser/.local/bin:/home/claudeuser/.bun/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
ExecStart=/bin/bash -c 'tmux has-session -t claude-telegram 2>/dev/null && tmux kill-session -t claude-telegram; tmux new-session -d -s claude-telegram "/home/claudeuser/start-claude.sh"'
ExecStop=/usr/bin/tmux kill-session -t claude-telegram

[Install]
WantedBy=multi-user.target
